# maintenance/tori-tier-mover.timer
[Unit]
Description=Tori Tier Mover Timer - Runs nightly at 01:00 CST
Documentation=https://github.com/tori/maintenance
Requires=tori-tier-mover.service

[Timer]
# Run daily at 01:00 CST (Central Standard Time)
# CST is UTC-6, so 01:00 CST = 07:00 UTC
OnCalendar=*-*-* 07:00:00
Persistent=true

# Add some randomization to prevent thundering herd
RandomizedDelaySec=300

# If the system was powered down, run the job when it starts up
Persistent=true

[Install]
WantedBy=timers.target
